# Server/mini.service
[Unit]
Description=Minimized server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/mono /home/pi/mini/Server.exe
WorkingDirectory=/home/pi/mini
StandardOutput=syslog
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target